# systemd/dell_xps_15_9570/disable_nvidia.service
[Unit]
Description=Disable NVIDIA card on boot.

[Service]
ExecStart=/home/cjuub/.disable_nvidia.sh

[Install]
WantedBy=multi-user.target
