# systemd/dell_xps_15_9570/lspci_after_suspend.service

[Unit]
Description="Fix issue after suspend"
After=suspend.target

[Service]
ExecStart=lspci

[Install]
WantedBy=suspend.target
