# targets/dentaku/disable_touchpad/services/disable_touchpad.service
[Unit]
After=suspend.target

[Service]
ExecStart=modprobe -r hid_multitouch

[Install]
WantedBy=suspend.target
